# allinone/supervisord.conf
; supervisord — runs all five services in the all-in-one container.
; Rendered config files live under /etc/ccobs/runtime/ (produced by entrypoint.sh).
; Persistent data lives under /data/{prometheus,loki,grafana}.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/ccobs/supervisord.log
pidfile=/run/ccobs/supervisord.pid
loglevel=info

[unix_http_server]
file=/run/ccobs/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/ccobs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; --- Loki (start first; collector depends on it) ---------------------------
[program:loki]
command=/usr/local/bin/loki -config.file=/etc/ccobs/runtime/loki-config.yaml
priority=10
autorestart=true
startsecs=5
stdout_logfile=/var/log/ccobs/loki.log
stderr_logfile=/var/log/ccobs/loki.log

; --- Prometheus ------------------------------------------------------------
[program:prometheus]
command=/usr/local/bin/prometheus --config.file=/etc/ccobs/runtime/prometheus.yml --storage.tsdb.path=/data/prometheus --storage.tsdb.retention.time=%(ENV_PROM_RETENTION)s --web.listen-address=127.0.0.1:9090
priority=20
autorestart=true
startsecs=5
stdout_logfile=/var/log/ccobs/prometheus.log
stderr_logfile=/var/log/ccobs/prometheus.log

; --- OTel Collector (receives OTLP, exports to Prometheus + Loki) ----------
[program:otel-collector]
command=/usr/local/bin/otelcol-contrib --config=/etc/ccobs/runtime/otel-config.yaml
priority=30
autorestart=true
startsecs=5
stdout_logfile=/var/log/ccobs/otel-collector.log
stderr_logfile=/var/log/ccobs/otel-collector.log

; --- Grafana ---------------------------------------------------------------
[program:grafana]
command=/usr/share/grafana/bin/grafana server --homepath=/usr/share/grafana --config=/etc/ccobs/runtime/grafana.ini
priority=40
autorestart=true
startsecs=5
stdout_logfile=/var/log/ccobs/grafana.log
stderr_logfile=/var/log/ccobs/grafana.log

; --- nginx (ingest gate + grafana proxy) -----------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=50
autorestart=true
startsecs=3
stdout_logfile=/var/log/ccobs/nginx.log
stderr_logfile=/var/log/ccobs/nginx.log
